# scripts/apex-dashboard.service
[Unit]
Description=Apex Hunter V14 Dashboard
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/apexBot
# Dashboard usually uses the same venv
ExecStart=/home/ubuntu/apexBot/venv/bin/streamlit run dashboard/overview.py --server.port 8505 --server.address 0.0.0.0
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=apex-dashboard

[Install]
WantedBy=multi-user.target
